# config/systemd/cortex.target
[Unit]
Description=Cortex Voice Assistant
After=multi-user.target

[Install]
WantedBy=multi-user.target
